# pi/rc.local
# Earthstar
/home/pi/tesseract-control-software/pi/effectbox.sh &
/home/pi/tesseract-control-software/pi/spidev-driver.sh &
